# cmd/assets/logmcp.service
[Unit]
Description=LogMCP — AI Log Access Server
Documentation=https://github.com/kleist-dev/logmcp
After=network.target

[Service]
ExecStart=/usr/local/bin/logmcp serve
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
User=logmcp
Group=logmcp
NoNewPrivileges=true
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/logmcp
ReadOnlyPaths=/var/log
RestrictNamespaces=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
